# dx-cluster-bridge.service
[Unit]
Description=DX Cluster WebSocket Bridge
After=network.target

[Service]
Type=simple
User=radio
# Edit --cluster-host and --callsign before deploying
ExecStart=/usr/local/bin/dx-cluster-bridge.py \
    --cluster-host dxcluster.example.com \
    --cluster-port 7300 \
    --callsign YOURCALL-6 \
    --ws-port 9373 \
    --max-age 30
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
